# tortuga/supervisord.conf
[supervisord]
nodaemon=true

[program:sabnzbd]
command=/usr/bin/sabnzbdplus --config-file /tortuga --server :8080

[program:sickbeard]
command=python2 /opt/sickbeard/SickBeard.py --datadir=/tortuga